# dss_voip_test/rootfs/etc/services.d/stdin/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# SDeSalve HomeAssistant Addon: VoIP Notifier for HomeAssistant
# Take down the S6 supervision tree when the addon fails
# ==============================================================================

s6-svscanctl -t /var/run/s6/services
